# api-test-site/openrc/api-test-site
#!/sbin/openrc-run

name="api-test-site"
description="Temporary API test site"
command="/usr/bin/node"
command_args="/opt/api-test-site/server.js"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
directory="/opt/api-test-site"

: ${API_TEST_PORT:=28884}

export PORT="${API_TEST_PORT}"

depend() {
  need net
}
